# StarsPinesV2/provision/systemd/starsandpines-backup.timer
[Unit]
Description=Stars & Pines Backup Timer
After=network.target

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target